# configs/gentoo/etc/local.d/hostap.start
#!/bin/sh

ip link set wlan0 up
ip addr add 192.168.1.1/24 dev wlan0

busybox udhcpd -S /etc/udhcpd.conf
hostapd -P /run/hostapd.pid -B /etc/hostapd/hostapd.conf
